# sherbert
#!/sbin/openrc-run

name="sherbert"
description="sherbet http server"

supervisor="supervise-daemon"

directory="/opt/sherbert/"

command="/usr/local/bin/sherbert"
command_args="host ${directory}/site/ --theme=nord --addr=0.0.0.0:8001"
command_user="sherbert"


start() {
	ebegin "Starting ${name}"
	supervise-daemon "${RC_SVCNAME}" \
                         --stdout /var/log/${RC_SVCNAME}.out.log \
	                 --stderr /var/log/${RC_SVCNAME}.err.log \
                         --user "${command_user}" \
			 -- "${RC_SVCNAME}" -- ${command_args}
	eend $?
}

stop() {
	ebegin "Stopping ${name}"
	supervise-daemon "${RC_SVCNAME}" --stop
	eend $?
}

depend() {
    need net
    after dns logger netmount
}
